# src/ipfs/crontab
1 * * * * /app/login_sync_s3.sh >> /var/log/cron.log 2>&1
# Don't remove the empty line at the end of this file. It is required to run the cron job > scheduler.txt
